# roles/os-keystone/templates/ubuntu/uwsgi.socket.j2
[Unit]
Description=uWSGI service

[Socket]
ListenStream=/run/uwsgi/%I.sock
SocketMode=0600
SocketUser={{ http_user }}

[Install]
WantedBy=sockets.target
